# crontab
* * * * * sh -c "java -cp /app/application.jar ar.edu.utn.frba.dds.pepita.EjemploJob"
