# crontab
* * * * * sh -c "java -cp /usr/src/target/application.jar ar.edu.utn.frba.dds.pepita.EjemploJob"
